# lamp/conf/mysql_supervisor.conf
[program:nginx]
command=/usr/local/mysql/bin/mysqld --user=mysql --datadir=/usr/local/mysql/data/ --init-file=/usr/local/mysql/init-file.sql
#daemon=true
autorestart=true

